# aeroponics.service
[Unit]
Description=Bill's Aeroponics Smart Watering System
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/plant-waterer
ExecStart=/usr/bin/python3 /home/pi/plant-waterer/app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONUNBUFFERED=1

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/pi/plant-waterer

[Install]
WantedBy=multi-user.target 